# deploy/social-datanode-live/auto-deploy/social-datanode-live-autodeploy.cron
*/3 * * * * $HOME/.local/bin/social-datanode-live-autodeploy >> $HOME/.local/share/social-datanode-live-autodeploy/cron.log 2>&1
